# contrib/elegoo-web.service
[Unit]
Description=Elegoo CC2 Web Service
Documentation=https://github.com/runnane/elegoo-web
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=elegoo-web
Group=elegoo-web
WorkingDirectory=/opt/elegoo-web
EnvironmentFile=/opt/elegoo-web/.env
ExecStart=/usr/bin/node --import tsx src/server/index.ts
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=elegoo-web

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/elegoo-web

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
